# files/system/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnBootSec=5m
OnUnitInactiveSec=1h
RandomizedDelaySec=5m
